# recipes/rc/rc.local
#!/bin/bash
#
# /etc/rc.local: local startup script
#

# End of file
